# template/tamplet.service
[Unit]
Description = Keyboard LED Remapping

[Service]
workingDirectory = /
ExecStart = %%TARGET%%

[Install]
WantedBy = multi-user.target